# python/paperclip.ini
[program:paperclip]
process_name=%(program_name)s_%(process_num)02d
command=node --import ./server/node_modules/tsx/dist/loader.mjs server/dist/index.js
directory=/app
user=paperclip
autostart=true
autorestart=true
priority=1
pidfile=/tmp/paperclip.pid
redirect_stderr=true
stdout_logfile=/var/www/logs/process-paperclip.log
environment=HOME="/paperclip",USER="paperclip",NODE_ENV="production",HOST="0.0.0.0",PORT="3100",PATH="/usr/local/share/pnpm:/usr/local/share/pnpm/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
